# scripts/service/polaris.timer.tpl
[Unit]
Description=Create virtual folder system

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target